# dot_config/autostart/xrandr_lum.desktop
[Desktop Entry]
Type=Application
Name=xrand
Exec=$HOME/.config/autostart/xrandr_lum.sh
StartupNotify=false
Terminal=false
